# rootfs/etc/services.d/openbox/run
#!/usr/bin/with-contenv sh

set -u # Treat unset variables as an error.

# During its startup, openbox tries to change its workdir to '$HOME'.  Make it
# happy by setting good one.  No data should be written there.
export HOME=/config

echo "[$(basename "$(pwd)")] starting..."

exec 2>&1 /bin/s6-notifyoncheck -d s6-setuidgid $USER_ID:$GROUP_ID /usr/bin/openbox-session

# vim: set ft=sh :
